# scripts/systemd/copilot-workspace.service.example
# ~/.config/systemd/user/fictional-scientists-copilot.service にコピーして使う
[Unit]
Description=copilot workspace
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/development/python-template-for-ai-assistant
ExecStart=/usr/bin/env bash -lc './scripts/copilot-compose.sh up -d copilot'
ExecStop=/usr/bin/env bash -lc './scripts/copilot-compose.sh down'
TimeoutStartSec=0

[Install]
WantedBy=default.target
